# crontabs/root
#!/bin/sh
0 * * * * 	/bin/sh /spotprice-mqtt/spotprice-mqtt.sh
