# crontabs/root
#!/bin/sh
*/30 * * * * 	/bin/sh /spotprice-mqtt/spotprice-mqtt.sh
